# tcb.service
[Unit]
Description=TCB Core Clipboard Storage Engine
After=graphical-session.target

[Service]
Type=simple

# Clean up old sockets/PIDs left behind by crashes
ExecStartPre=/usr/bin/bash -c 'rm -f /tmp/tcb.sock /tmp/tcb.pid'

# Start ONLY the core daemon backend
ExecStart=/usr/local/bin/tcb-daemon

Restart=on-failure
RestartSec=1s

[Install]
WantedBy=graphical-session.target