# tcb.service
[Unit]
Description=TCB Clipboard Sync Daemon
After=graphical-session.target

[Service]
Type=simple

# Force clean up old socket or PID files if they were left behind by a bad crash
ExecStartPre=/usr/bin/bash -c 'rm -f /tmp/tcb.sock /tmp/tcb.pid'

# Run the watch script
ExecStart=/usr/local/bin/tcb-watch

# Kill ALL processes in the control group on stop (daemon, watch, clipnotify)
KillMode=control-group

Restart=on-failure
# Add a 1-second delay before restarting to prevent the "repeated too quickly" loop
RestartSec=1s

[Install]
WantedBy=graphical-session.target